# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; =====================================================================
; TRADING API SERVERS
; =====================================================================

[program:longterm-api]
command=/Users/ashokkumar/Desktop/alg-discovery/.venv/bin/python -m uvicorn longterm_server:app --host 0.0.0.0 --port 8001 --log-level info
directory=/Users/ashokkumar/Desktop/alg-discovery/api
user=ashokkumar
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/longterm-api.log
stderr_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/longterm-api-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=PYTHONPATH="/Users/ashokkumar/Desktop/alg-discovery"

[program:swing-api]
command=/Users/ashokkumar/Desktop/alg-discovery/.venv/bin/python -m uvicorn swing_server:app --host 0.0.0.0 --port 8002 --log-level info
directory=/Users/ashokkumar/Desktop/alg-discovery/api
user=ashokkumar
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/swing-api.log
stderr_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/swing-api-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=PYTHONPATH="/Users/ashokkumar/Desktop/alg-discovery"

[program:shortterm-api]
command=/Users/ashokkumar/Desktop/alg-discovery/.venv/bin/python -m uvicorn shortterm_server:app --host 0.0.0.0 --port 8003 --log-level info
directory=/Users/ashokkumar/Desktop/alg-discovery/api
user=ashokkumar
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/shortterm-api.log
stderr_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/shortterm-api-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=PYTHONPATH="/Users/ashokkumar/Desktop/alg-discovery"

; ---------------------------------------------------------------------
; INTRADAY BUY API (port 8004)
; ---------------------------------------------------------------------

[program:intraday-buy-api]
command=/Users/ashokkumar/Desktop/alg-discovery/.venv/bin/python -m uvicorn alg_discovery.api.v1.routes.intraday_buy_server:app --host 0.0.0.0 --port 8004 --log-level info
directory=/Users/ashokkumar/Desktop/alg-discovery
user=ashokkumar
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/intraday-buy-api.log
stderr_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/intraday-buy-api-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=PYTHONPATH="/Users/ashokkumar/Desktop/alg-discovery/src:/Users/ashokkumar/Desktop/alg-discovery",INTRADAY_CHARTINK_ONLY="true"

; ---------------------------------------------------------------------
; INTRADAY SELL API (port 8005)
; ---------------------------------------------------------------------

[program:intraday-sell-api]
command=/Users/ashokkumar/Desktop/alg-discovery/.venv/bin/python -m uvicorn intraday_sell_server:app --host 0.0.0.0 --port 8005 --log-level info
directory=/Users/ashokkumar/Desktop/alg-discovery/api
user=ashokkumar
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/intraday-sell-api.log
stderr_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/intraday-sell-api-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=PYTHONPATH="/Users/ashokkumar/Desktop/alg-discovery",INTRADAY_CHARTINK_ONLY="true"

; =====================================================================
; CRON MANAGER
; =====================================================================

[program:trading-cron]
command=/Users/ashokkumar/Desktop/alg-discovery/.venv/bin/python cron_manager_starter.py
directory=/Users/ashokkumar/Desktop/alg-discovery
user=ashokkumar
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/trading-cron.log
stderr_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/trading-cron-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=PYTHONPATH="/Users/ashokkumar/Desktop/alg-discovery"

; =====================================================================
; CACHE REFRESH SERVICE (Every 1 Hour)
; =====================================================================

[program:cache-refresh]
command=/Users/ashokkumar/Desktop/alg-discovery/.venv/bin/python cache_refresh_service.py
directory=/Users/ashokkumar/Desktop/alg-discovery
user=ashokkumar
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/cache-refresh.log
stderr_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/cache-refresh-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=PYTHONPATH="/Users/ashokkumar/Desktop/alg-discovery"

; =====================================================================
; INTERNET CONNECTIVITY MONITOR
; =====================================================================

[program:internet-monitor]
command=/Users/ashokkumar/Desktop/alg-discovery/.venv/bin/python internet_monitor.py
directory=/Users/ashokkumar/Desktop/alg-discovery
user=ashokkumar
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/internet-monitor.log
stderr_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/internet-monitor-error.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile_backups=3
environment=PYTHONPATH="/Users/ashokkumar/Desktop/alg-discovery"

[program:misc-api]
command=/Users/ashokkumar/Desktop/alg-discovery/.venv/bin/python -m uvicorn misc_server:app --host 0.0.0.0 --port 8006 --log-level info
directory=/Users/ashokkumar/Desktop/alg-discovery/api
user=ashokkumar
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/misc-api.log
stderr_logfile=/Users/ashokkumar/Desktop/alg-discovery/logs/misc-api-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=PYTHONPATH="/Users/ashokkumar/Desktop/alg-discovery"

; =====================================================================
; GROUP CONFIGURATION
; =====================================================================

[group:trading-apis]
programs=longterm-api,swing-api,shortterm-api,intraday-buy-api,intraday-sell-api,misc-api
priority=999

[group:trading-services]
programs=trading-cron,cache-refresh,internet-monitor
priority=998 